# supervisor/conf.d/php-nowplaying.conf
[program:php-nowplaying]
command = /usr/bin/php /var/pheme/www/bin/console pheme:sync:nowplaying
user = pheme
priority = 600
numprocs = 1
autostart = true
autorestart = true
redirect_stderr = true

stopasgroup = true
killasgroup = true

# Logfiles
stdout_logfile = /var/pheme/logs/service_php-nowplaying.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups = 5
